# templates/default/services/systemd/nomad.service.erb
[Unit]
Description=Nomad Agent Service
After=syslog.target network.target

[Service]
User=<%= @user %>
Group=<%= @group %>
LimitNOFILE=<%= @nofile %>
ExecStart=/usr/bin/nomad agent -config <%= @config %>

[Install]
WantedBy=multi-user.target
